# mydns-sendip.service
[Unit]
Description=MyDNS.JP send IP address

[Service]
Type=oneshot
ExecStart=/opt/mydns-sendip-httpbasic/sendip.rb
